# sensord/data/openrc/sensord.init
#!/sbin/openrc-run

command=/usr/bin/sensord
command_args=
command_user=sensord:sensord

command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

output_logger='logger -p daemon.info -t sensord.stdout'
error_logger='logger -p daemon.err -t sensord.stderr'

depend() {
	need dbus
	use syslog
}
